# supervisor-dogecoin.conf
[supervisord]
nodaemon=true

[program:dogecoind]
command=su -c "dogecoind -conf=/home/crypto/dogecoin.conf" crypto
autorestart=true
